# supervisor.conf
[program:setup]
user=root
command=/opt/setup
redirect_stderr=true
autorestart=false

[program:apache]
user=root
command=/opt/start_apache
redirect_stderr=true
autorestart=true
